# ps5-idle-timeout.service
[Unit]
Description=DualSense Idle Timeout Monitor
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/ps5-idle-timeout
Restart=on-failure

[Install]
WantedBy=default.target